# target/linux/ipq806x/base-files/etc/hotplug.d/iface/11-eth-init
#!/bin/sh

. /lib/functions.sh
. /lib/functions/system.sh

board=$(board_name)

case $board in
        netgear,wac510)
                if [ $ACTION == ifup ]; then
			# Toggle eth LED pin to activate eth[0,1] LEDs when the interface comes up
                        echo 0 > /sys/class/leds/wac510\:green\:eth/brightness
                        echo 1 > /sys/class/leds/wac510\:green\:eth/brightness
                fi
                ;;
esac
;;

