# target/linux/ipq806x/base-files/etc/hotplug.d/net/01-wlan-mac
#!/bin/sh

. /lib/functions.sh
. /lib/functions/system.sh

board=$(board_name)

if [ "$ACTION" == "add" ]; then
        case "$board" in
        "netgear,wac510")
                if [ "$INTERFACE" == "wlan0" ] || [ "$INTERFACE" == "wlan1" ]; then
                        offset=$((${INTERFACE: -1} + 2))
                        new_mac=$(macaddr_add $(cat /sys/class/net/eth0/address) +$offset)
                        old_mac=$(cat /sys/class/net/$INTERFACE/address)
                        if [ ! "$old_mac" == "$new_mac" ]; then
				cur_state=$(cat /sys/class/net/$INTERFACE/operstate)
				if [ "$cur_state" == "up" ]; then
                                	ifconfig $INTERFACE down
                                	ifconfig $INTERFACE hw ether $new_mac
                                	ifconfig $INTERFACE up
				else
					ifconfig $INTERFACE hw ether $new_mac
				fi
                        fi
                fi
                ;;
        esac
fi

